# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:postgresql]
command=/usr/lib/postgresql/17/bin/postgres -D /data/postgresql
user=postgres
autorestart=true
priority=1

[program:streamlit]
command=streamlit run /app/src/streamlit_app.py --server.port=8501 --server.address=0.0.0.0
directory=/app
autorestart=true
priority=2